# transformations/rootfs/etc/services.d/rsyslogd/run
#!/usr/bin/execlineb -P
/usr/sbin/rsyslogd -f /etc/rsyslog.conf -n
